# Automation1/Exercises/05b-PracticalTemplatingandLookups/05b-ex5/templates/hello-service.service.j2
[Unit]
Description=Hello Service
After=network.target

[Service]
Type=exec
User=root
WorkingDirectory=/opt/hello-service
ExecStart=/usr/bin/python3 /opt/hello-service/app.py
Restart=on-failure

[Install]
WantedBy=multi-user.target